# templates/job.conf.erb
# <%= title %>
<% if description != ""%>description "<%= description %>"<% end %>

start on runlevel [35]
stop on runlevel [016]

# daemon or not?
<% if fork %>expect fork<% end %>
respawn

<% if use_script then%>script<% else %>exec<% end %> <%= cmd %>
